# debian/mobile-usb-networking.mobile-usb-network-setup.service
[Unit]
Description=Setup USB networking
Requires=NetworkManager.service
After=NetworkManager.service
ConditionFirstBoot=yes

[Service]
Type=oneshot
ExecStart=/usr/sbin/mobile-usb-network-setup

[Install]
WantedBy=multi-user.target
