# debian/mobile-usb-networking.prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
    # Remove USB gadget network connection
    nmcli connection delete RNDIS
    nmcli connection delete ECM

    # Disable automatic start of USB gadget
    deb-systemd-helper disable mobile-usb-gadget
fi

#DEBHELPER#
